# deploy/systemd/admin-dashboard.service
[Unit]
Description=Bot Admin Dashboard
Wants=redis-server.service
After=network.target redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/bot
EnvironmentFile=/opt/bot/.env
ExecStart=/opt/bot/venv/bin/uvicorn admin_dashboard:app --host 127.0.0.1 --port 8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
