# deploy/systemd/bot.service
[Unit]
Description=Telegram Bot
Wants=redis-server.service
After=network.target redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/bot
EnvironmentFile=/opt/bot/.env
User=bot
Group=bot
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/bot
ExecStart=/opt/bot/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
